# ds4.target
[Unit]
Description=DS4
RefuseManualStart=no
StopWhenUnneeded=no