# ds4hook.service
[Unit]
Description=ds4hook

[Service]
Type=simple
ExecStart=/usr/local/bin/ds4hook.sh

[Install]
WantedBy=ds4.target