# ds4leds.service
[Unit]
Description=ds4hook

[Service]
Type=simple
ExecStart=/usr/local/bin/ds4leds.sh

[Install]
WantedBy=ds4.target